# deploy/lanthon-synth.service
# =============================================================================
# lanthon-synth.service — LANTH0N 5YNTH SuperCollider audio engine
# =============================================================================
# Installed to /etc/systemd/system/ by deploy/setup.sh
# Starts JACK + sclang with src/main.scd on boot.

[Unit]
Description=LANTH0N 5YNTH — SuperCollider Live Performance Rig
After=network.target sound.target syslog.target lanth0n-cpugov.service
Wants=sound.target lanth0n-cpugov.service

[Service]
Type=simple
User=%USER%
Group=audio
WorkingDirectory=%PROJECT_DIR%
Environment="SC_JACK_DEFAULT_OUTPUT_PORT=system:playback_1"

# Start JACK, then sclang
ExecStartPre=/bin/sleep 3
ExecStartPre=/usr/local/bin/lanthon-jack-start.sh
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/sclang -D %PROJECT_DIR%/src/main.scd

Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

StandardOutput=append:/var/log/lanth0n/synth.log
StandardError=append:/var/log/lanth0n/synth.log
SyslogIdentifier=lanth0n-synth

LimitRTPRIO=95
LimitMEMLOCK=infinity
LimitNICE=-19
ReadWritePaths=/var/log/lanth0n %PROJECT_DIR%/config %PROJECT_DIR%/setlists %PROJECT_DIR%/media %PROJECT_DIR%/samples

[Install]
WantedBy=multi-user.target
